# s6-overlay/s6-rc.d/poller/finish
#!/command/with-contenv sh
echo "[s6] Poller exited (code $1 signal $2). Restarting in 1s..."
