# s6-overlay/s6-rc.d/poller/run
#!/command/with-contenv sh
export PYTHONPATH=/app
echo "[s6] Starting Nextcloud Talk poller..."
exec python3 -u /app/poller/main.py
